# crontab
MAILTO=""
37 17 * * * jsub -o dbsync.log -j y dbsync.job
