# crontab
MAILTO=""
37 17 * * * jsub -l release=trusty dbsync.job
